# web-files/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
user = root
autostart = true

[program:php-fpm]
command = /usr/sbin/php-fpm -F -R
user = root
autostart = true

[program:memcache]
command = /usr/bin/memcached -u root -l 0.0.0.0 -m 64
user = root
autostart = true

[program:gearmand]
command = /usr/sbin/gearmand --verbose INFO
user = root
autostart = true

[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --pid-file=/var/run/mysqld/mysqld.pid
autostart=true
autorestart=true
user=root

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface